# service/otel/run
#!/bin/sh

exec otelcol-contrib --config /etc/otelcol-contrib/config.yaml
